# lass2mysql.service
[Unit]
Description=lass2mysql
[Service]
ExecStart=/htdocs/PMS5003TAirQuality/lass2mysql.sh
[Install]
WantedBy=multi-user.target
